# scripts/buildtest.crontab
# Crontab for daikonbuildtest user
# Runs on beet.

# To initiate this crontab, run the crontab command on this file.
# Henceforth, this file will nightly be updated and its commands
# executed.

USER=daikonbuildtest
MAILTO=daikon@pag.lcs.mit.edu
#MAILTO=mharder@lcs.mit.edu

01 1 * * * cd $HOME/research/invariants/scripts; cvs -Q up
02 1 * * * crontab $HOME/research/invariants/scripts/buildtest.crontab
55 1 * * * $HOME/research/invariants/scripts/buildtest-remove-old
# 0 2 * * * $HOME/research/invariants/scripts/buildtest
## dfectests are failing as of 1/2004; the problem appears to be use of gcc
## 3.  The tests should be made to use gcc 2.95.3 again, I believe.
0 2 * * * $HOME/research/invariants/scripts/buildtest --skip_dfec
#44 13 * * * $HOME/research/invariants/scripts/buildtest
