# scripts/buildtest.crontab
# Crontab for daikonbuildtest user
# Runs on parsnip.

# To initiate this crontab, run the crontab command on this file.
# Henceforth, this file will nightly be updated and its commands
# executed.

USER=daikonuser
MAILTO=daikon-developers@csail.mit.edu
CRONIC=/afs/csail.mit.edu/u/d/daikonuser/research/invariants/plume-lib/cronic
AFS=/afs/csail.mit.edu/u/d/daikonuser/bin/afs-cron-wrapper

01 1 * * * cd $HOME/research/invariants/scripts; $AFS cvs -Q up
02 1 * * * cd $HOME/research/invariants/plume-lib; $AFS $CRONIC hg pull -u
03 1 * * * $AFS crontab $HOME/research/invariants/scripts/buildtest.crontab
55 1 * * * $AFS $HOME/research/invariants/scripts/buildtest-remove-old
0 2 * * * $AFS $HOME/research/invariants/scripts/buildtest
#02 15 * * * $AFS $HOME/research/invariants/scripts/buildtest

