# scripts/buildtest.crontab
# Crontab for daikonbuildtest user

USER=daikonbuildtest
MAILTO=daikon@geyer.lcs.mit.edu
#MAILTO=mharder@lcs.mit.edu

01 0 * * * crontab $HOME/research/invariants/scripts/buildtest.crontab
50 1 * * * cd $HOME/research/invariants/scripts; cvs -Q up
55 1 * * * $HOME/research/invariants/scripts/buildtest-remove-old
0 2 * * * $HOME/research/invariants/scripts/buildtest
#44 13 * * * $HOME/research/invariants/scripts/buildtest
